# .envs/dev/php/php-worker.conf
[program:cron]
autorestart=false
command=cron -f

[program:php-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/app/artisan queue:work redis --sleep=3 --tries=1 --timeout=180
autostart=true
autorestart=true
user=php-fpm
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/app/storage/logs/worker.log
stdout_logfile_backups=0
stdout_logfile_maxbytes=0

[program:php-monitor]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/app/artisan queue:work redis --queue=monitor --sleep=3 --tries=1 --timeout=180
autostart=true
autorestart=true
user=php-fpm
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/app/storage/logs/monitor.log
stdout_logfile_backups=0
stdout_logfile_maxbytes=0
